# systemd/dh-renew@.timer
[Unit]
Description=DH param update for %I

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=multi-user.target
